# supervisor/vin-decoder.ini
[program:tornado-8000]
command=/var/www/vin-decoder/h5/server.py
directory=/var/www/vin-decoder/h5
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:tornado-10090]
command=/var/www/vin-decoder/tornado/app.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-31f]
command=/var/www/vin-decoder/proxy/confirm.py 31f
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-kuaidaili]
command=/var/www/vin-decoder/proxy/confirm.py kuaidaili
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-xicidaili]
command=/var/www/vin-decoder/proxy/confirm.py xicidaili
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-66ip]
command=/var/www/vin-decoder/proxy/confirm.py 66ip
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-httpsdaili]
command=/var/www/vin-decoder/proxy/confirm.py httpsdaili
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:confirm-proxy-nianshao]
command=/var/www/vin-decoder/proxy/confirm.py nianshao
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info




[program:spider-proxy-31f]
command=/var/www/vin-decoder/spider/proxy/cn_31f.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-66ip]
command=/var/www/vin-decoder/spider/proxy/cn_66ip.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-kuaidaili]
command=/var/www/vin-decoder/spider/proxy/com_kuaidaili.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-xicidaili]
command=/var/www/vin-decoder/spider/proxy/com_xicidaili.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info

[program:spider-proxy-httpsdaili]
command=/var/www/vin-decoder/spider/proxy/com_httpsdaili.py
directory=/var/www/vin-decoder
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/vin-decoder.log
loglevel=info
